# agents/localai02/localai-trace-shipper/localai-trace-shipper.service
[Unit]
Description=LocalAI Backend Trace Shipper (localai02 → Tempo)
After=network.target

[Service]
Type=simple
User=ckennedy
ExecStart=/usr/bin/python3 /opt/observability-agent/localai-trace-shipper/localai_trace_shipper.py
Environment=LOCALAI_URL=http://localhost:8080
Environment=OTLP_ENDPOINT=http://10.25.25.80:4318
Environment=HOST_NAME=localai02
Environment=SERVICE_NAME=localai-localai02
Environment=POLL_INTERVAL=15
Environment=LOCALAI_API_KEY=sk-23d4af75e0ea3e6fc8187f9c5fb0d5ca
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
